# build/website/jameswhale.service
[Unit]
Description=uWSGI instance to serve jameswhale
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/var/www/website_redux
Environment="PATH=/var/www/website_redux/personal_website-env/bin"
ExecStart=/var/www/website_redux/personal_website-env/bin/uwsgi --ini build/website/jameswhale.ini

[Install]
WantedBy=multi-user.target
