# build/website/jameswhale.service
[Unit]
Description=uWSGI instance to serve jameswhale
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/var/www/personal_website-main
Environment="PATH=/var/www/personal_website-main/personal_website-env/bin"
ExecStart=/var/www/personal_website-main/personal_website-env/bin/uwsgi --ini build/website/james-whale.ini

[Install]
WantedBy=multi-user.target
